# dot_config/systemd/user/autorestic.timer
[Unit]
Description=Autorestic Backup Timer

[Timer]
OnCalendar=20:*:*
Persistent=true

[Install]
WantedBy=timers.target
