# dot_config/systemd/user/clean-downloads.timer
[Unit]
Description=Downloads Cleanup Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
